# cron/plists/com.mission-control.scout-weekly-content-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CLAUDE_MODEL</key>
		<string>sonnet</string>
		<key>HOME</key>
		<string>/Users/josephinewood</string>
		<key>PATH</key>
		<string>/Users/josephinewood/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>VAULT_DIR</key>
		<string>/Users/josephinewood/fern-vault</string>
	</dict>
	<key>Label</key>
	<string>com.mission-control.scout-weekly-content-review</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/josephinewood/mission-control/cron/scripts/scout-weekly-content-review.sh</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/josephinewood/Library/Logs/mission-control/scout-weekly-content-review.log</string>
	<key>StandardOutPath</key>
	<string>/Users/josephinewood/Library/Logs/mission-control/scout-weekly-content-review.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>17</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>5</integer>
	</dict>
	<key>ThrottleInterval</key>
	<integer>60</integer>
	<key>WorkingDirectory</key>
	<string>/Users/josephinewood/fern-vault</string>
</dict>
</plist>
