# cron/plists/com.mission-control.scout-derek-job-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mission Control — Scout Derek Job Scan

  Schedule: Wednesday 8:00am CT (13:00 UTC, Weekday=3)
  Cron equivalent: 0 13 * * 3

  Installation: run cron/install.sh — or manually:
    ln -sf "$(pwd)/plists/com.mission-control.scout-derek-job-scan.plist" ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mission-control.scout-derek-job-scan.plist

  Logs: ~/Library/Logs/mission-control/scout-derek-job-scan.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mission-control.scout-derek-job-scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/josephinewood/mission-control/cron/scripts/scout-derek-job-scan.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/josephinewood/fern-vault</string>

  <!-- Wednesday 8:00am CT = Wednesday 13:00 UTC (Weekday 3 = Wednesday) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>3</integer>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/josephinewood</string>
    <key>PATH</key>
    <string>/Users/josephinewood/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>VAULT_DIR</key>
    <string>/Users/josephinewood/fern-vault</string>
    <key>CLAUDE_MODEL</key>
    <string>sonnet</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/josephinewood/Library/Logs/mission-control/scout-derek-job-scan.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/josephinewood/Library/Logs/mission-control/scout-derek-job-scan.log</string>
</dict>
</plist>
